# contrib/init/mashhadbitd.conf
description "Mashhadbit Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env MASHHADBITD_BIN="/usr/bin/mashhadbitd"
env MASHHADBITD_USER="mashhadbit"
env MASHHADBITD_GROUP="mashhadbit"
env MASHHADBITD_PIDDIR="/var/run/mashhadbitd"
# upstart can't handle variables constructed with other variables
env MASHHADBITD_PIDFILE="/var/run/mashhadbitd/mashhadbitd.pid"
env MASHHADBITD_CONFIGFILE="/etc/mashhadbit/mashhadbit.conf"
env MASHHADBITD_DATADIR="/var/lib/mashhadbitd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # mashhadbitd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$MASHHADBITD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run mashhadbitd."
        echo "The setting must appear in $MASHHADBITD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Mashhadbit Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$MASHHADBITD_PIDDIR"
    chmod 0755 "$MASHHADBITD_PIDDIR"
    chown $MASHHADBITD_USER:$MASHHADBITD_GROUP "$MASHHADBITD_PIDDIR"
    chown $MASHHADBITD_USER:$MASHHADBITD_GROUP "$MASHHADBITD_CONFIGFILE"
    chmod 0660 "$MASHHADBITD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$MASHHADBITD_PIDFILE" \
    --chuid $MASHHADBITD_USER:$MASHHADBITD_GROUP \
    --exec "$MASHHADBITD_BIN" \
    -- \
    -pid="$MASHHADBITD_PIDFILE" \
    -conf="$MASHHADBITD_CONFIGFILE" \
    -datadir="$MASHHADBITD_DATADIR" \
    -disablewallet \
    -daemon

